# systemd/plasmashell-oom/.config/systemd/user/plasmashell-oom.service
[Unit]
Description=Mirai plasmashell memory monitor
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=%h/.local/bin/plasmashell-oom.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target